# lib/to-sort/os-integration/coordinator-recv-db.conf
start on networking
stop on runlevel [S016]
respawn
exec /usr/local/libexec/coordinator-recv-db
